# systemd/telora-daemon.service
[Unit]
Description=Telora Audio Daemon (Whisper)
After=network.target pipewire.service
Requires=pipewire.service

[Service]
Type=simple
ExecStart=/usr/bin/telora-daemon
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
